# Debian.files/DEBIAN/postrm
#!/bin/sh
# postrm script for mailscanner
#
set -e

# unregister the startup script
#update-rc.d -f mailscanner remove

exit 0