# Debian.files/DEBIAN/postrm
#!/bin/sh
# postrm script for mailscanner
#
set -e

# unregister the old startup script
update-rc.d -f MailScanner remove

exit 0