# RPM.files/common/processing_messages_alert.cron
#!/bin/bash

[ -x /usr/sbin/processing_messages_alert ] || exit 0
exec /usr/sbin/processing_messages_alert > /dev/null 2>&1
exit 0
